# deploy/systemd/blinky-dashboard.service
[Unit]
Description=Blinky LED Dashboard v3
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/blinky-led
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/blinky-led/.venv/bin/blinky serve --config /home/pi/blinky-led/config.yaml --env-file /home/pi/blinky-led/.env --log-level INFO --host 0.0.0.0 --port 8080
Restart=always
RestartSec=2
KillSignal=SIGINT
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
